# ui/qt-headunit/systemd/ui-log-collector.timer
[Unit]
Description=Periodic UI bringup log collection

[Timer]
OnBootSec=25s
OnUnitActiveSec=90s
Unit=ui-log-collector.service

[Install]
WantedBy=timers.target
