# ui/qt-headunit/systemd/ui-log-collector.timer
[Unit]
Description=Periodic UI bringup log collection

[Timer]
OnBootSec=12s
OnUnitActiveSec=60s
Persistent=true
AccuracySec=10s
RandomizedDelaySec=0
Unit=ui-log-collector.service

[Install]
WantedBy=timers.target
